# roles/openclaw/templates/s1-dashboard.service.j2
[Unit]
Description=S1 Docker Dashboard (TTY1)
After=docker.service network-online.target
Wants=network-online.target
# Restart if docker comes up late
After=docker.service

[Service]
Type=simple
User={{ openclaw_install_user }}
Group={{ openclaw_install_user }}
# Run on TTY1 so it shows on the physical console at boot
TTYPath=/dev/tty1
TTYReset=yes
TTYVHangup=yes
StandardInput=tty
StandardOutput=tty
StandardError=journal

Environment=TERM=xterm-256color
ExecStart=/usr/bin/python3 {{ openclaw_workspace }}/tools/docker-dashboard.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
